# supervisord.conf
[supervisord]
nodaemon=true

[program:app]
command=/usr/local/bin/main
autostart=true
autorestart=true
stderr_logfile=/var/log/app.err.log
stdout_logfile=/var/log/app.out.log
